# config/swapster.conf
[program:<название>]
command=<путь до гуникорна> <путь до wsgi файла>.wsgi:application - c <путь до конфига гуникорна>
directory=<путь до корневой директории проекта>
user=<имя пользователя>
autorestart=true
redirect_stderr=true
stdout_logfile = <имя до файла с логами>